# templates/systemd-redis.service.j2
[Unit]
Description=Advanced key-value store
After=network.target
Documentation=http://redis.io/documentation, man:redis-server(1)

[Service]
Type=notify
ExecStart=/usr/bin/redis-server {{ item['redis_confdir'] | default(redis_confdir) }}{{ item.redis_name | default(redis_name) }}.conf --supervised systemd --daemonize no
PIDFile=/run/redis/{{ item['redis_name'] | default(redis_name) }}.pid
TimeoutStopSec=0
Restart=always
User=redis
Group=redis

[Install]
WantedBy=multi-user.target